# actual/actual.container
[Unit]
Description=Actual Budget - Personal Finance Manager
Documentation=https://actualbudget.org/docs/
After=actual-data-volume.service
Requires=actual-data-volume.service
StartLimitBurst=5
StartLimitIntervalSec=300

[Container]
ContainerName=actual-server

Image=docker.io/actualbudget/actual-server:latest
AutoUpdate=registry

# Localhost only - external access via Caddy
PublishPort=127.0.0.1:5006:5006

Volume=actual-data.volume:/data

# Health check using wget
HealthCmd=wget -q --spider http://localhost:5006/ || exit 1
HealthInterval=60s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=30s
Notify=healthy
HealthOnFailure=kill

NoNewPrivileges=true
LogDriver=journald

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=default.target
