# immich/immich-server.container
[Unit]
Description=Immich Server - Photo and Video Management
Documentation=https://immich.app/docs
After=immich-network.service immich-redis.service immich-postgres.service immich-machine-learning.service immich-upload-volume.service
Requires=immich-network.service immich-redis.service immich-postgres.service immich-upload-volume.service
Wants=immich-machine-learning.service
StartLimitBurst=5
StartLimitIntervalSec=300

[Container]
ContainerName=immich-server
HostName=immich-server

Image=ghcr.io/immich-app/immich-server:v1.124.2
AutoUpdate=registry

Network=immich.network

# Localhost only - external access via Caddy
PublishPort=127.0.0.1:2283:2283

EnvironmentFile=/home/immich/.config/containers/systemd/immich.env

Volume=immich-upload.volume:/usr/src/app/upload:
Volume=/etc/localtime:/etc/localtime:ro

HealthCmd=wget -q --spider http://localhost:2283/api/server/ping || exit 1
HealthInterval=60s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s
Notify=healthy
HealthOnFailure=kill

NoNewPrivileges=true
LogDriver=journald

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=default.target
